# docs/examples/aehub-broker-credentials.service.example
# /etc/systemd/system/aehub-broker-credentials.service
#
# AE.HUB Broker Credentials Service
# 
# Infrastructure node that fetches MQTT broker configuration from Config Service
# and publishes it to ROS2 topic /aehub/mqtt/broker_config.
#
# Dependencies:
#   - ROS2 workspace built and installed
#   - Credentials/config file installed with the package:
#       share/aehub_broker_credentials/config/broker_credentials.yaml
#
# Usage:
#   sudo cp docs/examples/aehub-broker-credentials.service.example /etc/systemd/system/aehub-broker-credentials.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable aehub-broker-credentials.service
#   sudo systemctl start aehub-broker-credentials.service
#
# Security:
#   - Keep secrets OUT of git.
#   - Service runs as non-root user (adjust User/Group as needed)

[Unit]
Description=AE.HUB Broker Credentials (Config Service → ROS2)
Documentation=https://github.com/your-org/aehub
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=boris
Group=boris

# ROS2 setup and launch
ExecStart=/bin/bash -lc '\
  source /opt/ros/jazzy/setup.bash && \
  source /home/boris/ros2_ws/install/setup.bash && \
  ros2 launch aehub_broker_credentials broker_credentials.launch.py \
'

# Restart policy
Restart=always
RestartSec=3
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aehub-broker-credentials

[Install]
WantedBy=multi-user.target
