# scripts/aehub_nav2.service
[Unit]
Description=AEHub Nav2 stack (symovo_nav2.launch.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=boris
WorkingDirectory=/home/boris/ros2_ws

# Load runtime parameters (API keys etc.) outside of git.
EnvironmentFile=-/etc/aehub_nav2.env

# Sensible defaults if /etc/aehub_nav2.env is missing
Environment=ROS_DOMAIN_ID=0
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_cpp
Environment=ROBOT_ID=fahrdummy-01-local
Environment=CONFIG_SERVICE_URL=http://localhost:7900
Environment=CONFIG_SERVICE_API_KEY=CHANGE_ME
Environment=USE_SIM_TIME=false
Environment=AUTOSTART=true
Environment=TLS_VERIFY=false
Environment=SYMOVO_AUTO_ENABLE_DRIVE_MODE=true

ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash && source /home/boris/ros2_ws/install/setup.bash && exec ros2 launch aehub_navigation symovo_nav2.launch.py robot_id:=$ROBOT_ID config_service_url:=$CONFIG_SERVICE_URL config_service_api_key:=$CONFIG_SERVICE_API_KEY use_sim_time:=$USE_SIM_TIME autostart:=$AUTOSTART tls_verify:=$TLS_VERIFY symovo_auto_enable_drive_mode:=$SYMOVO_AUTO_ENABLE_DRIVE_MODE'

# Send SIGINT (like Ctrl+C) for clean shutdown.
KillSignal=SIGINT
TimeoutStopSec=25

# Clean up lock files after stop (safe to remove after process exit).
ExecStopPost=/bin/bash -lc 'rm -f /tmp/symovo_nav2_stack.lock /tmp/base_controller.lock /tmp/navigation_integrated_node.lock'

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
